# docker/root/etc/s6-overlay/s6-rc.d/svc-web/run
#!/command/with-contenv sh
# shellcheck shell=bash

cd /app/apps/web;
export LD_PRELOAD=/usr/local/lib/libjemalloc.so
exec node server.js;
