# docker/root/etc/s6-overlay/s6-rc.d/svc-workers/run
#!/command/with-contenv sh
# shellcheck shell=bash

cd /app/apps/workers;
export LD_PRELOAD=/usr/local/lib/libjemalloc.so
exec node dist/index.js
